# charging-threshold.path
[Path]
PathModified=/etc/charging-threshold

[Install]
WantedBy=multi-user.target